# 00-system/02-scripts/com.basket.onboarding.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>EnvironmentVariables</key>
    <dict>
        <key>GOOGLE_WORKSPACE_CLI_KEYRING_BACKEND</key>
        <string>file</string>
    </dict>
  <key>Label</key>
  <string>com.basket.onboarding</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/server-mini/do-better-workspace/00-system/02-scripts/basket-onboarding-sender.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>6</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
  </array>
  <key>StandardOutPath</key>
  <string>/Users/server-mini/do-better-workspace/20-operations/25-basket-ops-manual/onboarding-daily/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/server-mini/do-better-workspace/20-operations/25-basket-ops-manual/onboarding-daily/launchd.err.log</string>
</dict>
</plist>
