# 00-system/02-scripts/com.basket.ops-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Basket 운영 일일보고 봇 — launchd 상시가동
  설치:
    cp 00-system/02-scripts/com.basket.ops-bot.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.basket.ops-bot.plist
  (먼저 .env에 BASKET_BOT_TOKEN 등 채울 것)
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.basket.ops-bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/server-mini/do-better-workspace/20-operations/24-second-brain/.venv311/bin/python</string>
        <string>-u</string>
        <string>/Users/server-mini/do-better-workspace/00-system/02-scripts/basket-ops-bot.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/server-mini/do-better-workspace</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>GOOGLE_WORKSPACE_CLI_KEYRING_BACKEND</key>
        <string>file</string>
        <key>PATH</key>
        <string>/Users/server-mini/.npm-global/bin:/Users/server-mini/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/basket-ops-bot.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/basket-ops-bot.log</string>
</dict>
</plist>
