# 00-system/02-scripts/com.projectrent.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  프로젝트 포트폴리오 대시보드 서버 — launchd 상시 가동
  설치(맥미니에서):
    cp com.projectrent.dashboard.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.projectrent.dashboard.plist
  중지/재시작:
    launchctl unload ~/Library/LaunchAgents/com.projectrent.dashboard.plist
  ※ 경로의 사용자명(choi_ai)이 맥미니와 다르면 아래 두 경로를 수정.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.projectrent.dashboard</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/choi_ai/do-better-workspace/00-system/02-scripts/dashboard-server.py</string>
    <string>8770</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/choi_ai/do-better-workspace/00-system/02-scripts</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/pr-dashboard.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/pr-dashboard.err</string>
</dict>
</plist>
