# 00-system/02-scripts/com.projectrent.weekly-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.projectrent.weekly-report</string>
    <!-- ARISA 주간 업무 대시보드 — 매주 월요일 08:30, 지난주(월~일) 집계 → HTML + 텔레그램 요약.
         1회성 배치(KeepAlive 아님). gws CLI 위해 PATH에 ~/.npm-global/bin 포함. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/choi_ai/do-better-workspace/20-operations/24-second-brain/.venv311/bin/python</string>
        <string>-u</string>
        <string>/Users/choi_ai/do-better-workspace/00-system/02-scripts/weekly-report-aggregate.py</string>
        <string>--week</string>
        <string>last</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/choi_ai/do-better-workspace</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/choi_ai/.npm-global/bin:/Users/choi_ai/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/pr-weekly-report.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/pr-weekly-report.log</string>
</dict>
</plist>
